# rootfs/etc/services.d/3proxy/run
#!/command/with-contenv sh
set -e -u -o pipefail

while ! $(ip link | grep -q tun0); do sleep 1; done
echo "OpenVPN is running"


# start
3proxy /etc/3proxy/3proxy.cfg